# docker/crontab
# 每天凌晨 2:00 自动同步钉钉上级关系
# 环境变量 APP_INTERNAL_URL 和 CRON_SECRET 由 docker-compose 注入
0 2 * * * wget -q -O /dev/null --header="x-cron-secret: $CRON_SECRET" --post-data="" "$APP_INTERNAL_URL/api/admin/dingtalk/sync-manager-relations" >> /var/log/manager-sync.log 2>&1
